# rc.local
# Put your custom commands here that should be excuted once
# the system init finishes. By default this file does nothing
sh /root/amp-lab-http-rover/confignetwork.sh

gpsd -n /dev/ttyS0

sh /root/amp-lab-http-rover/collisionflash.sh
sh /root/amp-lab-http-rover/run_server.sh
 exit 0
