# rc.local
# Put your custom commands here that should be excuted once
# the system init finishes. By default this file does nothing
sh /root/amp-lab-http-rover-Joystick/confignetwork.sh
sh /root/amp-lab-http-rover-Joystick/collisionflash.sh
sh /root/amp-lab-http-rover-Joystick/run_server.sh
 exit 0
